# 50-x-resize.rules
ACTION=="change", KERNEL=="card1", SUBSYSTEM=="drm", RUN+="/usr/local/bin/x-resize"
